# appalachia-radio.service
[Unit]
Description=Appalachia Radio - Continuous Music Player
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/appalachia-radio
Environment="PATH=/home/pi/appalachia-radio/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/pi/appalachia-radio/venv/bin/python /home/pi/appalachia-radio/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

