# station/retrowaves-station.service
[Unit]
Description=Retrowaves Station (DJ + Playout Engine)
After=network.target sound.target retrowaves-tower.service
Wants=retrowaves-tower.service

[Service]
Type=simple
User=retrowaves
Group=retrowaves

# Your actual project root
WorkingDirectory=/opt/retrowaves

# Ensure imports resolve correctly - must be parent directory for station.* imports
Environment=PYTHONPATH=/opt/retrowaves

# Optional config
EnvironmentFile=-/etc/retrowaves/station.env

# Correct venv path - use station module's __main__.py entry point
ExecStart=/opt/retrowaves/venv/bin/python3 -m station

Restart=always
RestartSec=2

LimitNOFILE=65536
CPUSchedulingPolicy=rr
CPUSchedulingPriority=50

KillMode=mixed
TimeoutStopSec=300

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
