# tower/retrowaves-tower.service
[Unit]
Description=Retrowaves Tower (External Encoder + Stream Fanout)
After=network.target sound.target
# Tower starts independently; Station may optionally Want= it
Wants=network.target

[Service]
Type=simple
User=retrowaves
Group=retrowaves

# Create runtime directory for Unix socket with proper permissions
# This creates /run/retrowaves (or /var/run/retrowaves) owned by retrowaves:retrowaves
RuntimeDirectory=retrowaves
RuntimeDirectoryMode=0755

# Tower lives in its own directory
WorkingDirectory=/opt/retrowaves/tower

# Ensure imports resolve correctly
Environment="PYTHONPATH=/opt/retrowaves"

# Load environment variables from configuration file
# The '-' prefix means the file is optional (won't fail if missing)
EnvironmentFile=-/etc/retrowaves/tower.env

# Correct venv path (same as Station)
ExecStart=/opt/retrowaves/venv/bin/python3 -m tower

Restart=always
RestartSec=2

LimitNOFILE=65536
CPUSchedulingPolicy=rr
CPUSchedulingPriority=50

KillMode=mixed
TimeoutStopSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
